# boardcast_device.service
[Unit]
Description=Boardcast IoT Device Service
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/boardcast_device
EnvironmentFile=/home/dev/boardcast_device/.env
ExecStart=/home/dev/boardcast_device/venv/bin/python3 boardcast_device.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=boardcast_device

[Install]
WantedBy=multi-user.target