# boardcast_device.service
[Unit]
Description=Boardcast Device Agent
After=network.target

[Service]
Type=simple
User=rpi
WorkingDirectory=/home/rpi/boardcast_device
EnvironmentFile=/home/rpi/boardcast_device/.env
ExecStart=/home/rpi/boardcast_device/venv/bin/python3 boardcast_device.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
